# linuxservice/falconpitwitter.service
[Unit]
Description=Falcon Pi Twitter
After=network.target
Documentation=https://thealmostengineer.com/falconpitwitter

[Service]
Type=simple
Restart=always
ExecStart=/home/fpp/media/plugins/falconpitwitter/Almostengr.FalconPiTwitter
User=fpp

[Install]
WantedBy=multi-user.target